# crontab
# supercronic crontab — times are interpreted in the TZ set on the container
# (see docker-compose.yml). Default schedule below = hourly at :30,
# 09:30 through 21:30 in the container's local timezone.
#
# Format: minute hour day-of-month month day-of-week command
30 9-21 * * * python -m threads_poster
