# cron/threads_poster.cron
# Run the Threads poster every hour from 09:30 to 21:30 IST (04:00-16:00 UTC)
0 4-16 * * * cd /Users/devashishsharma/Documents/GitHub/12amstories/threads-api-posting && /usr/bin/env python3 -m threads_poster